# maint/crontab
#
#  Crontab for PETSc 2.0 nightly build
#
# This file is run every night.
# Run "crontab crontab" to install the nightly build;
# use "crontab -r" to remove it.
#
#  The timing between these runs IS not important.
#
1  0 * * * /home/bsmith/petsc/maint/buildremotetest
15 0 * * * /home/bsmith/petsc/maint/buildtest 1
15 0 * * * /home/bsmith/petsc/maint/buildtest 2
15 0 * * * /home/bsmith/petsc/maint/buildtest 3
15 0 * * * /home/bsmith/petsc/maint/buildtest 4
15 0 * * * /home/bsmith/petsc/maint/buildtest 5
15 0 * * * /home/bsmith/petsc/maint/buildtest 6
15 0 * * * /home/bsmith/petsc/maint/buildtest 7
45 8 * * * /usr/ucb/rsh eagle /home/bsmith/petsc/maint/checknightly
#
45 23 * * * /home/bsmith/petsc/maint/accesslogs/www/buildwwwlogs
