# maint/crontab
0 2 * * * /home/bsmith/petsc/maint/buildtest
