# maint/crontab
# $Id: crontab,v 1.153 2001/03/19 23:09:52 bsmith Exp bsmith $ 
#
#  Crontab for PETSc 2.0 nightly build
#
# This file is run every night on 'soot.mcs.anl.gov'.
# Run "crontab crontab" to install the nightly build;
# use "crontab -r" to remove it.
# 
#  The timing between these runs IS not important.
#
1  23 * * * /usr/ucb/rsh soot /home/bsmith/petsc/maint/buildremotetest
1  23 * * * /bin/rm -f /home/bsmith/petsc/maint/nightlylogs/*.log
1  23 * * * /home/bsmith/petsc/maint/removelibs
#
30  1 * * * /home/bsmith/petsc/maint/buildtest rs6000
40  1 * * * /home/bsmith/petsc/maint/buildtest solaris
50  1 * * * /home/bsmith/petsc/maint/buildtest IRIX
01  2 * * * /home/bsmith/petsc/maint/buildtest linux
10  2 * * * /home/bsmith/petsc/maint/buildtest linux_absoft
20  2 * * * /home/bsmith/petsc/maint/buildtest solaris_uni
30  2 * * * /home/bsmith/petsc/maint/buildtest IRIX64
40  2 * * * /home/bsmith/petsc/maint/buildtest rs6000_sp
#
#
# configure buids on various machines
#
50  2 * * * /home/bsmith/petsc/maint/buildconfigtest linux
#
#
# Check for and Zap unfinished processes
#
45  8 * * * /usr/ucb/rsh soot /home/bsmith/petsc/maint/checknightly
15  9 * * * /usr/ucb/rsh fire "/home/bsmith/petsc/maint/killnightly >& /home/bsmith/petsc/maint/nightlylogs/killall_log"
#
# Recover the logfiles for the nightly builds
#
30  9 * * * /home/bsmith/petsc/maint/copytestlog solaris
30  9 * * * /home/bsmith/petsc/maint/copytestlog linux
30  9 * * * /home/bsmith/petsc/maint/copytestlog linux_absoft
30  9 * * * /home/bsmith/petsc/maint/copytestlog rs6000
30  9 * * * /home/bsmith/petsc/maint/copytestlog solaris_uni
30  9 * * * /home/bsmith/petsc/maint/copytestlog IRIX
30  9 * * * /home/bsmith/petsc/maint/copytestlog IRIX64
30  9 * * * /home/bsmith/petsc/maint/copytestlog rs6000_sp
#
# Buildlinks on various machines
#
30 4 * * * /usr/ucb/rsh denali /home/bsmith/petsc/maint/buildlinks /sandbox/petsc/petsc
30 4 * * * /usr/ucb/rsh chiba /home/bsmith/petsc/maint/buildlinks /sandbox/petsc/petsc
# 
#
# Copy over www access log files for the past 24 hours
#
### 45 23 * * * /home/bsmith/petsc/maint/accesslogs/www/buildwwwlogs
#
# Backup the PETSc source dir
#
45 5 * * *  /usr/ucb/rsh soot /home/bsmith/petsc/maint/buildbackup
#
# Rebuild the fortran stubs - before the nightly builds are started.
#
01 22 * * * /usr/ucb/rsh fire "cd /home/bsmith/petsc;make allfortranstubs > /dev/null"
#
# Stash away a copy of petsc.tar.gz every tuesday
#
45 5 * * 2 /usr/ucb/rcp /home/petsc/petsc.tar.gz mcshsm:~/userbackups/
#
# Rebuild the documentation every night
#
15 6 * * * /usr/ucb/rsh terra "cd /home/bsmith/petsc;make allmanualpages LOC=/home/bsmith/petsc >& /home/bsmith/petsc/maint/nightlylogs/doc.log"
#
# send out reminder about Tuesday afternoon meeting
#
0 14 * * 1 /home/bsmith/petsc/maint/tuesdayreminder
15 14 * * 2 /home/bsmith/petsc/maint/tuesdayreminder
#
# check for bad PetscFunctionBegin
#
0 7 * * * /usr/ucb/rsh fire "cd /home/bsmith/petsc;make ACTION=checkbadPetscFunctionBegin tree >& /home/bsmith/petsc/maint/nightlylogs/petscbegin.log"
#
# check for missing RCS header in each file
#
0 7 * * * /usr/ucb/rsh fire "cd /home/bsmith/petsc;make ACTION=checkbadFileHeader tree >& /home/bsmith/petsc/maint/nightlylogs/fileheader.log"
#
# check for style errors in PETSc source
#
0 7 * * * /usr/ucb/rsh fire "cd /home/bsmith/petsc;make ACTION=checkbadSource tree >& /home/bsmith/petsc/maint/nightlylogs/source.log"
