# maint/crontab
#
#  Crontab for PETSc 2.0 nightly build
#
# This file is run every night.
# Run "crontab crontab" to install the nightly build;
# use "crontab -r" to remove it.
#
#  The timing between these runs IS not important.
#
1  23 * * * /home/bsmith/petsc/maint/buildremotetest
1  23 * * * /bin/rm /home/bsmith/petsc/maint/nightlylogs/*.log
1  23 * * * /home/bsmith/petsc/maint/removelibs
#
30  1 * * * /home/bsmith/petsc/maint/buildtest 1
30 23 * * * /home/bsmith/petsc/maint/buildtest 2
30  2 * * * /home/bsmith/petsc/maint/buildtest 3
30  1 * * * /home/bsmith/petsc/maint/buildtest 4
30  2 * * * /home/bsmith/petsc/maint/buildtest 5
30  2 * * * /home/bsmith/petsc/maint/buildtest 6
30  2 * * * /home/bsmith/petsc/maint/buildtest 8
30  2 * * * /home/bsmith/petsc/maint/buildtest 9
#
# This build is a second build on maverick and so is 
# started a bit late...
#
30  4 * * * /home/bsmith/petsc/maint/buildtest 7
#
#  Checks if all the nightly builds are complete and 
# sends mail to PETSc developers if there are hanging jobs
#
45  8 * * * /usr/ucb/rsh eagle /home/bsmith/petsc/maint/checknightly
#
# Zap unfinished processes on eagle
#
45 9 * * * /usr/ucb/rsh eagle zap tcsh petsc
#
45 23 * * * /home/bsmith/petsc/maint/accesslogs/www/buildwwwlogs
