# maint/crontab
#
#  Crontab for PETSc 2.0 nightly build
#
# This file is run every night.
# Run "crontab crontab" to install the nightly build;
# use "crontab -r" to remove it.
#
#  The timing between these runs IS important.
#
40 0 * * * /home/bsmith/petsc/maint/buildremotetest
15 1 * * * /home/bsmith/petsc/maint/buildtest bsmith
15 1 * * * /home/bsmith/petsc/maint/buildtest curfman
#15 1 * * * /home/bsmith/petsc/maint/buildtest petsc
#15 1 * * * /home/bsmith/petsc/maint/buildtest balay
#
# Makes the shadow directories.  Must be done before the library builds.
#
05 0 * * * /home/bsmith/petsc/maint/buildlinks /home/balay/petsc > /dev/null
05 0 * * * /home/bsmith/petsc/maint/buildlinks /home/curfman/petsc > /dev/null
05 0 * * * /home/bsmith/petsc/maint/buildlinks /home/petsc/petsc > /dev/null
